# fight_predictor.service.example
[Unit]
Description=Fight Predictor API
After=network.target postgresql.service

[Service]
# User and Group - create a dedicated user for your app
User=fight_predictor
Group=fight_predictor

# Working directory - path to your project
WorkingDirectory=/home/fight_predictor/fight_predictor

# Environment file - loads secrets from secure file
# Create this file: sudo nano /etc/fight_predictor/env
# Set permissions: sudo chmod 600 /etc/fight_predictor/env
EnvironmentFile=/etc/fight_predictor/env

# Python virtual environment path
Environment="PATH=/home/fight_predictor/.local/share/virtualenvs/fight_predictor-xxx/bin"

# Gunicorn command
ExecStart=/home/fight_predictor/.local/share/virtualenvs/fight_predictor-xxx/bin/gunicorn \
    --config gunicorn_config.py fight_predictor.wsgi:application

# Restart policy
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

